# app/rootfs/etc/supervisord.d/websockify.ini
[program:websockify]
command=/noVNC/websockify/websockify.py --web /noVNC/noVNC --unix-target=/tmp/.X11-unix/vnc.sock 6080
autorestart=true
user=novnc
environment=HOME="/home/novnc",USER="novnc"
directory=/home/novnc
