# app/rootfs/etc/xdg/autostart/ibus-daemon.desktop
[Desktop Entry]
Name=ibsu-daemon
Exec=/usr/bin/ibus-daemon -rxd
Type=Application
