# app/conf.d/xfce4.conf
[program:xfce4-session]
command=startxfce4
autorestart=true
user=novnc
environment=HOME="/home/novnc",USER="novnc"
directory=/home/novnc
